# rules/wazuh/local_rules.xml
<group name="windows,sysmon,powershell,">

  <rule id="100100" level="10">
    <if_sid>61603</if_sid>
    <field name="win.eventdata.image" type="pcre2">(?i)powershell\.exe$</field>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)-(?:EncodedCommand|enc|e)\s</field>
    <description>Suspicious PowerShell encoded command detected (T1059.001)</description>
    <mitre>
      <id>T1059.001</id>
    </mitre>
    <group>powershell_encoded,attack.execution,</group>
  </rule>

</group>

<group name="windows,sysmon,credential_access,">

<rule id="100101" level="12">
    <if_group>sysmon_eid10_detections</if_group>
    <field name="win.eventdata.targetImage" type="pcre2">(?i)lsass\.exe$</field>
    <description>Suspicious LSASS memory access — possible credential dumping (T1003.001)</description>
    <mitre>
      <id>T1003.001</id>
    </mitre>
    <group>credential_access,mimikatz,attack.credential_access,</group>
  </rule>
</group>


<group name="windows,sysmon,defense_evasion,">

  <rule id="100102" level="12">
    <if_sid>61603</if_sid>
    <field name="win.eventdata.image" type="pcre2">(?i)powershell\.exe$</field>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)Set-MpPreference.+Disable</field>
    <description>Windows Defender disabled via PowerShell (T1562.001)</description>
    <mitre>
      <id>T1562.001</id>
    </mitre>
    <group>defense_evasion,attack.defense_evasion,</group>
  </rule>

</group>

<group name="windows,sysmon,persistence,">

  <rule id="100103" level="10">
    <if_sid>61603</if_sid>
    <field name="win.eventdata.image" type="pcre2">(?i)schtasks\.exe$</field>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)/create</field>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)(powershell|cmd\.exe|wscript|mshta|rundll32)</field>
    <description>Suspicious scheduled task created with script interpreter (T1053.005)</description>
    <mitre>
      <id>T1053.005</id>
    </mitre>
    <group>persistence,attack.persistence,</group>
  </rule>

</group>

<group name="windows,defense_evasion,">

  <rule id="100104" level="14">
    <if_group>windows_logs</if_group>
    <field name="win.system.eventID">^1102$</field>
    <description>Windows Security event log cleared — possible evidence removal (T1070.001)</description>
    <mitre>
      <id>T1070.001</id>
    </mitre>
    <group>indicator_removal,attack.defense_evasion,</group>
  </rule>

  <rule id="100105" level="14">
    <if_group>windows_logs</if_group>
    <field name="win.system.eventID">^104$</field>
    <description>Windows System event log cleared — possible evidence removal (T1070.001)</description>
    <mitre>
      <id>T1070.001</id>
    </mitre>
    <group>indicator_removal,attack.defense_evasion,</group>
  </rule>

</group>
